# supervisord.conf
[supervisord]
nodaemon=true
user=root

[program:scraper]
command=python3 scraper/src/main.py background --interval 60
directory=/app
environment=PATH="/usr/local/bin:%(ENV_PATH)s",HOME="/root"
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/scraper.log

[program:engine]
command=python3 engine/src/main.py server
directory=/app
environment=PATH="/usr/local/bin:%(ENV_PATH)s",HOME="/root"
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/engine.log

[program:webserver]
command=node webserver/dist/index.js
directory=/app
environment=PATH="/usr/local/bin:%(ENV_PATH)s",HOME="/root"
autostart=true
autorestart=true
startretries=3
redirect_stderr=true
stdout_logfile=/var/log/webserver.log